# docker_scripts/athena.conf
[program:athena]
command=/athena-start start
autostart=true
autorestart=true
stderr_logfile=/var/log/athena.err.log
stdout_logfile=/var/log/athena.out.log
